# docker/supervisor_confs/msfs.ini
[program:msfs]
directory=/opt/im_server/%(program_name)s
command=/opt/im_server/%(program_name)s/%(program_name)s
autostart=true
autorestart=true
environment=LD_LIBRARY_PATH=/opt/im_server/lib